# etc/supervisord_sh.conf.tmpl
[program:pcapdb_uwsgi]
command=DESTDIR/bin/uwsgi --ini DESTDIR/etc/uwsgi.ini --uid capture --gid WWW_USER --logto DESTDIR/log/uwsgi.log
autorestart=true
user=root
stopsignal=KILL

