# etc/supervisord_sh.conf.tmpl
[program:pcapdb_uwsgi]
command=env uwsgi --ini DESTDIR/etc/uwsgi.ini --uid capture --gid WWW_USER
autorestart=true
user=root
stopsignal=KILL

